# configs/mysql/etc/services.d/mysql/run
#!/usr/bin/env bash
exec /usr/sbin/mysqld --user=mysql
