# configs/mysql/etc/services.d/mysql/run
#!/usr/bin/env bash
exec systemctl enable mysqld
exec systemctl start mysqld
