# configs/mariadb/etc/bk.init/op.conf
start on startup
task
exec /db_access.sh
